# templates/config/systemd/barrenero_nvidia_tuning.service.j2
[Unit]
Description=Nvidia tuning for Miner

[Install]
WantedBy=multi-user.target

[Service]
WorkingDirectory={{ common.lib }}
ExecStart=/usr/bin/python3 nvidia_tuning.py start
RemainAfterExit=yes